# roles/mounts/templates/nas.automount.j2
[Unit]
Description=Automount NAS NFS
After=network-online.target

[Automount]
Where={{ nas_mount_path }}

[Install]
WantedBy=multi-user.target

