# roles/mounts/templates/rclone-s3.service.j2
[Unit]
Description=rclone: Remote FUSE filesystem for cloud storage config
After=network-online.target
Wants=network-online.target
AssertPathIsDirectory={{ s3_mount_path }}

[Service]
Type=notify
ExecStart= \
  /usr/bin/rclone mount \
	--allow-other \
	--gid 1000 \
	--uid 1000 \
	--cache-dir "{{ ansible_env.HOME }}/.cache/rclone/s3" \
	--config "{{ ansible_env.HOME }}/.config/rclone/rclone.conf" \
	--vfs-cache-mode full \
	--vfs-cache-max-size 5G \
	--vfs-cache-max-age 8760h \
	--vfs-fast-fingerprint \
	--no-modtime \
	--rc \
	--rc-no-auth \
	--rc-addr :5576 \
	s3:cld {{ s3_mount_path }}
ExecStop=/bin/fusermount -uz {{ s3_mount_path }}
RestartSec=5
Restart=on-failure

[Install]
WantedBy=default.target
